# kreastrap/overlay/etc/jumpstart/services/getty-tty6.service
[Info]
Desc="This service starts up TTY6."
Type="service"

[Service]
exec="getty 115200,38400,9600 tty6 linux"