# src/jumpstart/services/test.service
[Service]
ExecPre="echo xd"
Exec="echo 'lol' > /lol"
ExecPost="echo lol!!!!"
Deps=""
Restart="yes"
