# src/kreastrap/overlay/etc/jumpstart/services/getty-tty4.service
[Info]
Desc="This service starts up TTY4."
Type="service"

[Service]
exec="getty 115200,38400,9600 tty4 linux"