# src/kreastrap/overlay/etc/jumpstart/services/getty-tty1.service
[Info]
Desc="This service starts up TTY1."
Type="service"

[Service]
exec="getty 115200,38400,9600 tty1 linux"